# supervisord.conf
[supervisord]
nodaemon=false
logfile=/app/data/supervisord.log
pidfile=/tmp/supervisord.pid

[program:uvicorn]
command=uvicorn backend.app.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/data/uvicorn.log
stderr_logfile=/app/data/uvicorn.err

[program:static]
command=python3 -m http.server 3000 --directory /app/frontend/dist
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/data/static.log
stderr_logfile=/app/data/static.err
